# data/org.gnome.Evince.service.in
[Unit]
Description=Evince document viewer
Documentation=man:evince(1)

[Service]
Type=dbus
BusName=@application_id@.Daemon
ExecStart=@libexecdir@/evinced
