# data/org.gnome.Evince.service.in
[Unit]
Description=Evince document viewer
Documentation=man:evince(1)

[Service]
Type=dbus
BusName=org.gnome.Evince.Daemon
ExecStart=@libexecdir@/evinced
